# supervisor/ini/orderapi.ini
[program:orderapi]
directory=./order/api
command=go build order.go -o build && go run order.go -f etc/order.yaml
redirect_stderr = true
stdout_logfile=./logs/orderApi.log